# Practica_1/Codigo/web_usuarioUGR/crontab
* * * * * /root/script_mantenimiento.sh
